# goip-bridge.service
[Unit]
Description=goip-bridge - GoIP SMS/USSD gateway
After=network-online.target nftables.service mariadb.service mysql.service
Wants=network-online.target

[Service]
Type=simple
User=goip-bridge
Group=goip-bridge
WorkingDirectory=/opt/goip-bridge
ExecStart=/opt/goip-bridge/goip-bridge -config /opt/goip-bridge/config.json
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
